# operations/systemd/openclaw-gateway-autoreload.path
[Unit]
Description=Watch openclaw npm package for updates and restart gateway

[Path]
PathModified=/usr/lib/node_modules/openclaw/package.json
Unit=openclaw-gateway-autoreload.service

[Install]
WantedBy=multi-user.target
