# operations/systemd/rook-dispatcher.timer
[Unit]
Description=Dispatch ready canonical tasks every 15 seconds

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
Unit=rook-dispatcher.service

[Install]
WantedBy=timers.target
