# operations/systemd/rook-dispatcher.timer
[Unit]
Description=Dispatch ready canonical tasks and inspect active hook claims every 2 minutes
Requires=rook-dispatcher.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=rook-dispatcher.service

[Install]
WantedBy=timers.target
