# operations/systemd/rook-dispatcher.service
[Unit]
Description=Rook Canonical Task Dispatcher
After=openclaw-gateway.service

[Service]
Type=oneshot
WorkingDirectory=/root/.openclaw/workspace
Environment=ROOK_DISPATCH_TIMEOUT_SECONDS=60
Environment=ROOK_STALE_CLAIM_MINUTES=20
Environment=ROOK_DISPATCH_MAX_ATTEMPTS=3
Environment=ROOK_DISPATCH_RETRY_DELAY_MS=5000
Environment=ROOK_NOTIFY_MAX_ATTEMPTS=3
Environment=ROOK_NOTIFY_RETRY_DELAY_MS=3000
TimeoutStartSec=120
KillMode=control-group
ExecStart=/usr/bin/node /root/.openclaw/workspace/operations/bin/task-dispatcher.mjs --limit 5 --source-channel system:dispatcher --dispatch-mode local
