# operations/systemd/rook-vm-ops-baseline.timer
[Unit]
Description=Rook VM Ops Baseline Report (daily 05:30)
Requires=rook-vm-ops-baseline.service

[Timer]
OnCalendar=*-*-* 05:30:00
Persistent=true
RandomizedDelaySec=5m
Unit=rook-vm-ops-baseline.service

[Install]
WantedBy=timers.target
