# operations/systemd/rook-dashboard-watchdog.service
[Unit]
Description=Rook Dashboard Watchdog
After=rook-dashboard.service

[Service]
Type=oneshot
ExecStart=/root/.openclaw/workspace/operations/bin/dashboard-watchdog.sh
